# HW4_file_system/3/pic_loader.service
[Unit]
Description=load pics from web

[Service]
Type=simple
ExecStart=/usr/bin/pic_loader.sh

[Install]
WantedBy=multi-user.target

